# z4s_daemon.service
[Unit]
Description=ZSpace Z4S - SATA Power, LED Control & Health Monitor
After=sysinit.target local-fs.target
DefaultDependencies=no

[Service]
Type=simple
ExecStart=/etc/z4s_daemon boot-monitor
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
